# scripts/com.buffalovfx.qtwatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.buffalovfx.qtwatcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Unversioned python3 symlink inside Shotgun.app: survives the app's
             bundled-Python version bumps (was python3.11, now python3.13). -->
        <string>/Applications/Shotgun.app/Contents/Resources/Python3/bin/python3</string>
        <string>/Volumes/atv-post-lucid3/atv-buffalo-s03/buffalo_vfx/buffalo_flow_config_alts/BUF_Mac_watcher/scripts/qt_watcher.py</string>
    </array>

    <!-- sgtk is not on the bundled Python's default path; point PYTHONPATH at
         the Toolkit core's python dir so `import sgtk` works under launchd. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Volumes/atv-post-lucid3/atv-buffalo-s03/buffalo_vfx/buffalo_flow_config/install/core/python</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Volumes/atv-post-lucid3/atv-buffalo-s03/buffalo_vfx/buffalo_flow_config/logs/qt_watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/atv-post-lucid3/atv-buffalo-s03/buffalo_vfx/buffalo_flow_config/logs/qt_watcher_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Volumes/atv-post-lucid3/atv-buffalo-s03/buffalo_vfx/buffalo_flow_config_alts/BUF_Mac_watcher/scripts</string>
</dict>
</plist>
